# Linux/meta-scholar-risc-v/recipes-connectivity/ip-config/files/ip-ethernet.service
[Unit]
Description=Configuration IP Ethernet
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/ip addr add 192.168.7.2/24 dev end0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
